# build/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=error
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
user=root
childlogdir=/var/log/

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:sls-live-server]
priority=50
autostart=true
autorestart=unexpected
redirect_stderr=true
user=nobody
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/var/log/sls-live-server_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/sls-live-server_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
command=/usr/local/bin/logprefix /usr/local/bin/sls -c /etc/sls/sls.conf

[program:srtla_rec]
priority=150
autostart=true
autorestart=unexpected
redirect_stderr=true
user=nobody
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/var/log/srtla_rec_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/srtla_rec_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
command=/bin/bash -c "sleep 3 && /usr/local/bin/logprefix /usr/local/bin/srtla_rec 5000 localhost 8282"